# weather-dashboard/supervisor.conf
[program:weather-dashboard-backend]
command=%(ENV_HOME)s/Coding/Vibe Coding/python-requests-main/weather-dashboard/venv/bin/python app.py --no-debug
directory=%(ENV_HOME)s/Coding/Vibe Coding/python-requests-main/weather-dashboard/backend
autostart=true
autorestart=true
startretries=5
startsecs=5
redirect_stderr=true
stdout_logfile=%(ENV_HOME)s/Coding/Vibe Coding/python-requests-main/weather-dashboard/logs/backend.log
environment=PORT=5003,FLASK_APP=app.py,FLASK_ENV=production,FLASK_DEBUG=0

[program:weather-dashboard-frontend]
command=npm start
directory=%(ENV_HOME)s/Coding/Vibe Coding/python-requests-main/weather-dashboard/frontend
autostart=true
autorestart=true
startretries=5
startsecs=10
redirect_stderr=true
stdout_logfile=%(ENV_HOME)s/Coding/Vibe Coding/python-requests-main/weather-dashboard/logs/frontend.log
environment=PORT=3000,BROWSER=none,DISABLE_ESLINT_PLUGIN=true,TSC_COMPILE_ON_ERROR=true,SKIP_PREFLIGHT_CHECK=true