# weather-dashboard/backend/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=logs/supervisord.log
pidfile=supervisord.pid
childlogdir=logs/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:weather-dashboard-backend]
command=python3 main.py
directory=%(here)s
autostart=true
autorestart=true
startretries=3
stderr_logfile=logs/backend.err.log
stdout_logfile=logs/backend.out.log
environment=PORT=5001
stopwaitsecs=10
stopsignal=TERM
killasgroup=true
stopasgroup=true

[group:weather-dashboard]
programs=weather-dashboard-backend